# microswift.conf
[program:microswift]
command = /home/microswift/microswift/run.sh
user = microswift
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
